# deploy/systemd/trading-wa-bridge.service
[Unit]
Description=Trading Agent - WhatsApp Bridge (OpenAI-compat)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/eswar/claude/trading-agent
EnvironmentFile=/home/eswar/claude/trading-agent/.env
ExecStart=/home/eswar/claude/trading-agent/.venv/bin/python whatsapp_bridge.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
