# deploy/systemd/trading-tunnel-healthcheck.timer
[Unit]
Description=Trading Agent - Reverse SSH Tunnel Healthcheck (every 60s)

[Timer]
OnBootSec=2min
OnUnitActiveSec=60s
AccuracySec=5s
Unit=trading-tunnel-healthcheck.service

[Install]
WantedBy=timers.target
